# rohtembak-frontend.service
[Unit]
Description=RohTembak Frontend
After=network.target rohtembak-backend.service
Wants=rohtembak-backend.service

[Service]
Type=simple
User=root
WorkingDirectory=/app/frontend
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/npm run dev -- --host 0.0.0.0 --port 3000
Restart=always
RestartSec=10
StandardOutput=append:/var/log/rohtembak/frontend.log
StandardError=append:/var/log/rohtembak/frontend.log

[Install]
WantedBy=multi-user.target
